# root/etc/services.d/php/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid root service php7.4-fpm start
